# install/macos/tangram-xtgmnode.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>tangram-xtgmnode</string>
	<key>ServiceDescription</key>
	<string>tangram xtgmnode daemon</string>
	<key>Program</key>
	<string>/opt/tangram/xtgmnode/xtgmnode</string>
	<key>WorkingDirectory</key>
	<string>/opt/xtgm/xtgmnode</string>
	<key>UserName</key>
	<string>_tangram_xtgmnode</string>
	<key>GroupName</key>
	<string>tangram_xtgmnode</string>
	<key>RunAtLoad</key>
	<true/>
	<key>ProcessType</key>
	<string>Backgroup</string>
</dict>
</plist>
